# scripts/com.thoughtcatalog.transcript-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thoughtcatalog.transcript-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/chris/tc-transcript-bot/scripts/tc-transcript-bot.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TC_CHRIS_NUMBER</key>
        <string>+1XXXXXXXXXX</string>
        <key>TC_PROJECT_DIR</key>
        <string>/Users/chris/tc-transcript-bot</string>
        <key>CLAUDE_BIN</key>
        <string>/usr/local/bin/claude</string>
        <key>DEEPGRAM_API_KEY</key>
        <string>REPLACE_ME</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>REPLACE_ME</string>
        <key>INSTAGRAM_USER</key>
        <string>REPLACE_ME</string>
        <key>INSTAGRAM_PASS</key>
        <string>REPLACE_ME</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/chris/tc-transcript-bot/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/chris/tc-transcript-bot/logs/launchd-stderr.log</string>
</dict>
</plist>
